# cfg-templates/crontab
30 * * * * bash -lc /scripts/topology_submit
